# crontab
*/5 *    * * *   root    /home/user/bitcoin-wallet-finder/verify.sh
0   5    * * *   root    /home/user/bitcoin-wallet-finder/upload.py
